# debian/dsme-dbus.upstart
description "Connect dsme to dbus"
author "Semi Malinen"

start on started dbus
start on started dbus-actdead

console output

script
        /usr/sbin/waitdbus system
        /usr/sbin/dsmetool --start-dbus
end script
